# mac/LaunchAgents/com.dockglo.learn.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
 "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>           <string>com.dockglo.learn</string>
  <key>ProgramArguments</key><array><string>/usr/bin/python3</string><string>~/dev/dockglo/scripts/learn_usage.py</string></array>
  <key>StartInterval</key>   <integer>3600</integer>
  <key>RunAtLoad</key>       <true/>
  <key>StandardErrorPath</key><string>/tmp/dg_learn.err</string>
  <key>StandardOutPath</key> <string>/tmp/dg_learn.out</string>
</dict>
</plist>
